# files/verlihub.in
#!/bin/sh

# PROVIDE: verlihub
# REQUIRE: FILESYSTEMS defaultroute netwait resolv
# REQUIRE: mysql unbound
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr
: ${verlihub_enable="NO"}
name="verlihub"
rcvar="verlihub_enable"

vh="/usr/local/bin/vh"
vbpcname="verlihub"

load_rc_config $name

start_cmd="vh_start"
stop_cmd="vh_stop"
restart_cmd="vh_restart"

vh_start()
{
	echo "Starting Verlihub..."
	su -s root -c '/usr/local/bin/vh -r'
}

vh_restart()
{
	echo "Restarting Verlihub..."
	su -s root -c '/usr/local/bin/vh -q'
}

vh_stop()
{
	echo "Stoping Verlihub..."
	su -s root -c '/usr/local/bin/vh -s'
}

#extra_commands="startvnc"
run_rc_command "$1"
